# final/pwn/sanity-check/solve/initramfs/etc/inittab
::sysinit:/etc/init.d/rcS
::once:-sh -c 'cat /etc/motd; setuidgid 1000 sh; poweroff'
